# roles/mqtt/templates/mosquitto.systemd.service.j2
#
# {{ ansible_managed }}
#

{#
currently disabled due to missing WITH_SYSTEMD support in Portage, see
https://bugs.gentoo.org/878969
[Service]
Type=notify
#NotifyAccess=main#}