# roles/redis/templates/redis.systemd.service.j2
#
# {{ ansible_managed }}
#

[Service]
Type=notify
