# roles/rsync_mirror/templates/ebuild-repo-sync.systemd.timer.j2
#
# {{ ansible_managed }}
#

[Unit]
Description=Run Gentoo ebuild repository sync nightly

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
