# roles/rsync_mirror/templates/ebuild-repo-sync.systemd.service.j2
#
# {{ ansible_managed }}
#

[Unit]
Description=Sync all Gentoo ebuild repositories

[Service]
Type=oneshot
ExecStart=/usr/bin/emaint sync --allrepos
