# roles/homeassistant/templates/homeassistant-update.timer.j2
#
# {{ ansible_managed }}
#

[Unit]
Description=Weekly update for Home Assistant

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
